# scripts/hermes-updater.service
[Unit]
Description=Hermes Auto-Updater

[Service]
Type=oneshot
ExecStart=%h/hermes/scripts/auto_update.sh
